# template/TEMPLATE.path
[Unit] 
Description=Path unit for monitoring document closing
[Path]
PathChanged=/home/root/.local/share/remarkable/xochitl/UUID.metadata
[Install]
WantedBy=multi-user.target
