# twitter/twitter.service
[Unit]
Description=Weather service
[Service]
ExecStart=/home/root/.local/share/pdflets/twitter.sh
